# apps/desktop/Resources/LaunchAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- The watch-folder LaunchAgent. build.sh stages this into
     Contents/Library/LaunchAgents/<bundle-id>.watcher.plist and rewrites Label +
     AssociatedBundleIdentifiers per channel so the three installs never collide.
     Registered/unregistered at runtime via SMAppService.agent from the app's
     Settings. BundleProgram is relative to the app bundle, so it keeps working if
     the app is moved. -->
<dict>
	<key>Label</key>
	<string>com.syntaxlabtechnology.crisp.watcher</string>
	<key>BundleProgram</key>
	<string>Contents/MacOS/CrispWatcher</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>AssociatedBundleIdentifiers</key>
	<array>
		<string>com.syntaxlabtechnology.crisp</string>
	</array>
	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
